# etc/cron.d/sync_git_repos
*/2 * * * * www-data /var/www/bin/update-repos.py
